# gvsense-agent.service
[Unit]
Description=gVsense Calibration Agent
Documentation=https://github.com/gvsense/gvsense
After=network.target
Wants=network.target

[Service]
Type=simple
User=gvsense
Group=gvsense
WorkingDirectory=/opt/gvsense
ExecStart=/usr/bin/python3 /opt/gvsense/web_server.py
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=gvsense-agent

# Security settings
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/var/lib/gvsense /var/log/gvsense

# Resource limits
LimitNOFILE=65536
MemoryMax=512M

[Install]
WantedBy=multi-user.target